# roles/setup_check_my_ip_server/files/check_my_ip_super_config.conf
[program:check_my_ip]
command=export FLASK_APP=/home/ansible/check_my_ip_server.py && flask run --host=0.0.0.0
