# roles/setup_ansible/files/message_bridge.conf
[program:message_bridge]
command=/home/ansible/ansible/message_bridge/message_consumer.py
